# modules/install-consul/consul.path
[Path]
PathExists=/etc/sysconfig/consul.conf
PathChanged=/opt/consul/certs

[Install]
WantedBy=multi-user.target
