# modules/install-consul/consul.path
[Path]
PathExists=/etc/sysconfig/consul.conf
PathChanged=/opt/consul/certs
PathChanged=/opt/consul/config

[Install]
WantedBy=multi-user.target
