# modules/install-consul/consul.path
[Path]
PathExists=/etc/sysconfig/consul.conf

[Install]
WantedBy=multi-user.target
