# modules/install-consul/update-ca-cert.service
[Unit]
Description=Update cacert

[Service]
Type=oneshot
ExecStart=/usr/sbin/update-ca-certificates

[Install]
WantedBy=network.target
